# deploy/record-collection.service
[Unit]
Description=Record Collection API
After=network.target

[Service]
User=glyphic
WorkingDirectory=/home/glyphic/record-collection/backend
Environment=FLASK_APP=main.py
Environment=DISCOGS_USER=your_discogs_user
Environment=DISCOGS_TOKEN=your_discogs_token
ExecStart=/home/glyphic/record-collection/backend/venv/bin/python -m flask run --host=0.0.0.0 --port=5000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
